# curator.cron
*/2 *  * * * /usr/local/bin/curator --config /etc/curator/curator.yml /etc/curator/actions.yml
